# contrib/testing/jenkins/runit/log/run
#!/bin/bash

export LC_ALL=en_US.UTF-8
export LANG=en_US.UTF-8
export LANGUAGE=en_US.UTF-8
exec svlogd /var/log/jenkins
